# services.ini
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'
directory=/app

[program:cron]
command=/usr/sbin/crond -f
directory=/app
user=root
autostart=true
autorestart=true
